# ubuntu/mongodb/init_mongo.conf
[program:init-mongo]
command = /bin/bash -c 'sleep 10 && /usr/bin/mongo localhost:27017 </tmp/init.js'
autorestart = unexpected
